# which_server.service.sample
[Unit]
Description=Which server service
After=systemd-user-sessions.service
[Service]
Type=simple
ExecStart=/path/to/which/start.sh