# which_server.service
[Unit]
Description=Which server service
After=systemd-user-sessions.service
[Service]
Type=simple
ExecStart=/home/jwlarocque/Documents/which/start.sh